# linux/filehostapi.service
# systemd 服务配置 - FileHostAPI 后端
# 路径: /etc/systemd/system/filehostapi.service
# 最后更新: 2025-11-06

[Unit]
Description=FileHost API Service
After=network.target

[Service]
WorkingDirectory=/var/www/filehostapi
ExecStart=/usr/bin/dotnet /var/www/filehostapi/FileHostAPI.dll
Restart=always
RestartSec=5
SyslogIdentifier=filehostapi
User=root
Environment="ASPNETCORE_ENVIRONMENT=Production"
Environment="ASPNETCORE_URLS=http://127.0.0.1:5001"
Environment="DOTNET_PRINT_TELEMETRY_MESSAGE=false"

[Install]
WantedBy=multi-user.target
